# overlay/usr/local/etc/rc.d/aero_firstboot
#!/bin/sh
#
# Aero OS first-boot experience — runs once, then never again
# Shows the Aero logo overlay, shrinks it away, and opens the setup screen.

# PROVIDE: aero_firstboot
# REQUIRE: aero_shell
# KEYWORD: shutdown

. /etc/rc.subr

name="aero_firstboot"
rcvar="aero_firstboot_enable"
marker="/var/db/aero-firstboot-done"
command="/usr/local/bin/aero-firstboot"
start_cmd="aero_firstboot_start"
stop_cmd=":"

aero_firstboot_start()
{
    if [ -f "${marker}" ]; then
        return 0
    fi
    if [ -z "${WAYLAND_DISPLAY}" ]; then
        export WAYLAND_DISPLAY=wayland-0
    fi
    if [ -z "${XDG_RUNTIME_DIR}" ]; then
        export XDG_RUNTIME_DIR="/tmp/runtime-$(id -u)"
        mkdir -p "${XDG_RUNTIME_DIR}"
        chmod 700 "${XDG_RUNTIME_DIR}"
    fi
    /usr/sbin/daemon -f "${command}"
}

load_rc_config "${name}"
run_rc_command "$1"
